# conf/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777
username = root
password = 123

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon=true
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = 123
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
username = root
password = 123
prompt = mysupervisor

[program:gunicorn]
directory=/app
command=gunicorn -c conf/gunicorn_config.py ocr.wsgi
priority=998

[program:nginx]
command     = /usr/sbin/nginx -c /etc/nginx/nginx.conf
directory   = /tmp
user        = root
startsecs   = 3
autorestart = true
priority=999
startretries            = 100
stdout_logfile          = /linkedsee/logs/supervisord/nginx.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 10
stderr_logfile          = /linkedsee/logs/supervisord/nginx_err.log
